# deploy/salt/eaa/web/templates/supervisor.conf.jinja
; Supervisor configuration for EAA donation portal
;
; ***** DO NOT EDIT THIS FILE *****
;
; This file is managed by Salt via {{ source }}

{% set user_home = salt['user.info'](config.user_name).home %}

[group:donations]
programs=web,celery

[program:web]
command={{ config.venv_path }}/bin/gunicorn wsgi -c {{ config.directory }}/donation_portal/gunicorn.conf.py --enable-stdio-inheritance
directory={{ config.directory }}/donation_portal
user={{ config.user_name }}
environment=USER="{{ config.user_name }}",
    HOME="{{ user_home }}",
    WORKON_HOME="{{ config.venv_path }}",
    PATH="{{ config.venv_path }}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",
    PYTHONUNBUFFERED=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile={{ daemon.supervisord.childlogdir }}/donations-web.log

[program:celery]
command={{ config.venv_path }}/bin/celery -A donation_portal worker -B -l INFO
directory={{ config.directory }}
user={{ config.user_name }}
environment=USER="{{ config.user_name }}",
    HOME="{{ user_home }}",
    WORKON_HOME="{{ config.venv_path }}",
    PATH="{{ config.venv_path }}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile={{ daemon.supervisord.childlogdir }}/donations-celery.log

; vim: cc=80 tw=79 ts=4 sw=4 sts=4 et sr
